# original_src/files/hashcat-cloud.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/hashcat_start

[Install]
WantedBy=default.target
